# scripts/samkwang-backend.service
[Unit]
Description=SAMKWANG AI FastAPI backend
After=network.target

[Service]
Type=simple
ExecStart=/bin/bash -c "set -a; [ -f '/home/user2/SAMKWANG AI/.env' ] && . '/home/user2/SAMKWANG AI/.env'; set +a; cd '/home/user2/SAMKWANG AI/backend' && export PYTHONPATH=. && exec '/home/user2/SAMKWANG AI/backend/.venv/bin/uvicorn' app.main:app --host 127.0.0.1 --port 8000"
Restart=on-failure
RestartSec=5

# journalctl -u samkwang-backend -f
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
